# morrigan/.config/autostart/imwheel.desktop
[Desktop Entry]
Exec=/usr/bin/imwheel -d -b 45
Icon=dialog-scripts
Name=IMWheel
Type=Application
X-KDE-AutostartScript=true
